# commafeed/rootfs/etc/services.d/commafeed/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: CommaFeed
# Runs CommaFeed
# ==============================================================================

bashio::log.info 'Starting CommaFeed...'

# Run CommaFeed
java -Djava.net.preferIPv4Stack=true -jar commafeed.jar server config.yml
